# supervisor.unison.tpl.conf
[program:unison-$name]
command = unison $volume -socket $port
user=$user
environment=HOME="$homedir",USER="$user"
directory = $volume
redirect_stderr = true
